# debian/plymouth.plymouth-stop.upstart
start on (stopped rc RUNLEVEL=[2345] or starting rcS)
stop on stopped plymouth

pre-start script
    exec /bin/plymouth quit
end script
